# python/spacewalk/satellite_tools/spacewalk-diskcheck.service
[Unit]
Description=Check available disk space on vital filesystems

[Service]
Type=oneshot
KillMode=process
ExecStart=/usr/bin/spacewalk-diskcheck
