# containers/server-systemd-services/uyuni-server.service
# container-uyuni-server.service
# autogenerated by Podman 4.3.1
# Tue Feb 28 17:20:52 CET 2023

[Unit]
Description=Uyuni server image container service
Wants=network.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=-/etc/sysconfig/uyuni-server-systemd-services
Restart=on-failure
ExecStartPre=/bin/rm \
    -f %t/uyuni-server.pid %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --conmon-pidfile %t/uyuni-server.pid \
    --cidfile=%t/%n.ctr-id \
    --cgroups=no-conmon \
    --rm \
    --sdnotify=conmon \
    -d \
    --replace \
    --tmpfs /run \
    -p 443:443 \
    -p 80:80 \
    -p 4505:4505 \
    -p 4506:4506 \
    -p 69:69 \
    -p 25151:25151 \
    -p 5432:5432 \
    -p 9100:9100 \
    -p 9187:9187 \
    -p 9800:9800 \
    -v cgroup:/sys/fs/cgroup:rw \
    -v var-lib-cobbler:/var/lib/cobbler \
    -v var-pgsql:/var/lib/pgsql \
    -v var-cache:/var/cache \
    -v var-spacewalk:/var/spacewalk \
    -v var-log:/var/log \
    -v srv-salt:/srv/salt \
    -v srv-www-pub:/srv/www/htdocs/pub \
    -v srv-www-cobbler:/srv/www/cobbler \
    -v srv-www-osimages:/srv/www/os-images \
    -v srv-tftpboot:/srv/tftpboot \
    -v srv-formulametadata:/srv/formula_metadata \
    -v srv-pillar:/srv/pillar \
    -v srv-susemanager:/srv/susemanager \
    -v srv-spacewalk:/srv/spacewalk \
    -v root:/root \
    -v etc-apache2:/etc/apache2 \
    -v etc-rhn:/etc/rhn \
    -v etc-systemd:/etc/systemd/system/multi-user.target.wants \
    -v etc-salt:/etc/salt \
    -v etc-tomcat:/etc/tomcat \
    -v etc-cobbler:/etc/cobbler \
    -v etc-sysconfig:/etc/sysconfig \
    -v etc-tls:/etc/pki/tls \
    -v ca-cert:/etc/pki/trust/anchors/ \
    --env-host \
    --hostname ${UYUNI_FQDN} \
    $EXTRA_POD_ARGS \
    --name uyuni-server ${NAMESPACE}/server:${TAG}
ExecStop=/usr/bin/podman stop \
    --ignore -t 10 \
    --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
    -f \
    --ignore -t 10 \
    --cidfile=%t/%n.ctr-id

PIDFile=%t/uyuni-server.pid
TimeoutStopSec=180
TimeoutStartSec=900
Type=forking

[Install]
WantedBy=multi-user.target default.target
