# proxy/proxy/saltproxy/saltproxy.service
[Unit]
Description=Salt ZeroMQ Proxy
After=syslog.target

[Service]
Type=simple
ExecStart=/usr/sbin/saltproxy

[Install]
WantedBy=multi-user.target
