# susemanager-utils/testing/docker/scripts/cobbler/supervisord/conf.d/dhcpd.conf
[program:dhcpd]
command=/usr/sbin/dhcpd -f -cf /etc/dhcpd.conf -user dhcpd --no-pid
autostart=false
autorestart=false
redirect_stderr=true
