# backend/satellite_tools/spacewalk-diskcheck.timer
[Unit]
Description=Check available disk space on vital filesystems

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

